# install/pem353.service
[Unit]
Description=PEM353 Modbus Gateway Service
After=network.target

[Service]
Type=simple
User=pem353
Group=pem353
SupplementaryGroups=dialout

WorkingDirectory=/opt/pem353

# Ensure log directory exists with correct permissions before starting
ExecStartPre=/bin/mkdir -p /var/log/pem353
ExecStartPre=/bin/chown pem353:pem353 /var/log/pem353
ExecStartPre=/bin/chmod 755 /var/log/pem353

ExecStart=/opt/pem353/bin/pem353

Restart=on-failure
RestartSec=5

# Environment
Environment=LD_LIBRARY_PATH=/opt/pem353/lib

# Security: Allow binding to privileged port 502 without root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pem353

[Install]
WantedBy=multi-user.target
